# bin/startup/debian/asynchronous-publishing
#! /bin/sh

### BEGIN INIT INFO
# Provides:          eZ Publish asynchronous publishing daemon
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: starts instance of eZ Publish asynchronous publishing daemon
# Description:       starts instance of eZ Publish asynchronous publishing daemon using start-stop-daemon
### END INIT INFO

### Configuration

# user to run as. Your webserver user, usually
RUN_AS=www-data

# startup args
DAEMON_OPTS=" -n"

### End configuration

STARTUP_SCRIPT_PATH=`dirname $(readlink -f $0)`
APP_PATH=${STARTUP_SCRIPT_PATH%/bin/startup/*}
DAEMON=$APP_PATH/bin/php/publishingqueueprocessor.php
NAME=`basename $0`
DESC="eZ Publish asynchronous publishing daemon"
PID_FILE=$APP_PATH/var/run/asynchronous-publishing.pid

. /lib/lsb/init-functions

# Test existence and executable status of executable
test -x $DAEMON || exit 0

UGID=$(getent passwd $RUN_AS | cut -f 3,4 -d:) || true

set -e

d_start() {
        if [ -z "$UGID" ]; then
                log_failure_msg "user \"$RUN_AS\" does not exist" "$NAME"
                return 1
        fi
        start-stop-daemon -d $APP_PATH -c $RUN_AS --start --background --pidfile $PID_FILE --exec $DAEMON -- $DAEMON_OPTS
        return $?
}

d_stop() {
        start-stop-daemon --stop --pidfile $PID_FILE

        if [ -f $PID_FILE ]; then rm -f $PID_FILE; fi
        return $?
}

case "$1" in
  start)
        log_daemon_msg "Starting $DESC" "$NAME"
        d_start
        log_end_msg $?
        ;;
  stop)
        log_daemon_msg "Stopping $DESC" "$NAME"
        d_stop
        log_end_msg $?
        ;;

  restart|force-reload)
        log_daemon_msg "Restarting $DESC" "$NAME"
        d_stop
        sleep 1
        d_start
        log_end_msg $?
        ;;
  *)
        N=/etc/init.d/$NAME
        echo "Usage: $N {start|stop|restart|force-reload}" >&2
        exit 1
        ;;
esac

exit 0